# cronjobs/cloudflare.cron
*/5 * * * * root /usr/bin/bash /workspaces/anonaddy-setup/scripts/cloudflare_update.sh
