# cronjobs/ipblacklist.cron
0 3 * * * root /usr/bin/bash /workspaces/anonaddy-setup/scripts/ip_blacklist_update.sh
